# nuk3gh0st.conf
description "Upstart nuk3gh0st service"
author      "Lil Peep"
start on filesystem or runlevel [2345]
stop on shutdown
exec /usr/share/nuk3gh0st/nuk3gh0st.startup start
